# planon.service
[Unit]
Description=Planon

[Service]
ExecStart=/path/to/planon

[Install]
WantedBy=multi-user.target
